# dev/build/zh/os/etc/profile.d/rust.proxy.sh
#!/usr/bin/env bash

export RUSTUP_UPDATE_ROOT="https://rsproxy.cn/rustup"
export RUSTUP_DIST_SERVER="https://rsproxy.cn"
